# src/busybox/etc/init.d/rc.local
#!/bin/sh

#echo "local service: telnetd start"
#telnetd

#echo "Samba start"
#smbd -D

